# services/supervisor/cron.d/schedule
* * * * * /usr/local/bin/php /app/laravel/artisan schedule:run >> /var/log/cron.log 2>&1
